# recipes-seco-synthetic-brain/seco-synthetic-brain/files/seco-synthetic-brain.service
[Unit]
Description=SECO-Synthetic-Brain

[Service]
Environment=NODE_ENV=production
ExecStart=/bin/su root -c 'node /usr/share/seco-synthetic-brain/index.js'
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
