# recipes-seco-synthetic-brain/seco-synthetic-brain/files/seco-synthetic-brain-watcher.path
[Path]
PathModified=/var/lib/seco-iot/config.json

[Install]
WantedBy=multi-user.target
