# recipes-seco-synthetic-brain/seco-synthetic-brain/files/seco-synthetic-brain-watcher.path
[Path]
PathModified=/var/seco-sb/config.json

[Install]
WantedBy=multi-user.target
