# recipes-seco-synthetic-brain/seco-synthetic-brain/files/seco-synthetic-brain-watcher.path
[Path]
PathModified=/var/lib/seco-iot/agent/config.json
PathModified=/var/lib/seco-iot/agent/actuators.json
PathModified=/var/lib/seco-iot/agent/sensors.json

[Install]
WantedBy=multi-user.target
